# roles/headplane_service/templates/headplane.service.j2
[Unit]
Description=Headplane Service
After=network.target
Requires=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=headplane
Group=headplane
WorkingDirectory={{ headplane_build_dir }}
Environment="PATH={{ headplane_node_path | dirname }}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="NVM_DIR=/home/headplane/.nvm"
Environment="HOME=/home/headplane"
ExecStart={{ headplane_node_path }} {{ headplane_build_dir }}/build/server/index.js
Restart=on-failure
RestartSec=5s

# Security restrictions for Node.js application
ProtectSystem=strict
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ReadWritePaths={{ headplane_build_dir }} {{ headplane_data_dir }} {{ headplane_log_dir }} {{ headplane_config_dir }} {{ headscale_socket_dir }}
ProtectHostname=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target